# pybinclock.service
[Unit]
Description=PyBinClock Binary Clock Display
After=network.target
Wants=network.target

[Service]
Type=simple
User=rmrfslashbin
Group=rmrfslashbin
WorkingDirectory=/home/rmrfslashbin/pybinclock
Environment=PATH=/usr/bin:/bin
ExecStart=/home/rmrfslashbin/.local/bin/uv run python -c 'from pybinclock.BinClockLEDs import BinClockLEDs; BinClockLEDs()'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target